# kiosk_edge/api-node/eisa-api.service
[Unit]
Description=E-ISA Kiosk Local API (Node.js / Fastify)
After=network.target

[Service]
Type=simple
User=eisa
Group=eisa
WorkingDirectory=/opt/eisa/api
EnvironmentFile=/etc/eisa/kiosk.env
ExecStart=/usr/bin/node /opt/eisa/api/src/index.js
Restart=always
RestartSec=3

# Sistemd hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
ReadWritePaths=/var/lib/eisa
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
